# packaging/systemd/slz-platform-backup.service.template
[Unit]
Description=SLZ Platform Database Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__ROOT_DIR__
ExecStart=__ROOT_DIR__/scripts/backup-platform-db.sh

[Install]
WantedBy=multi-user.target