# packaging/systemd/slz-platform-backup.timer
[Unit]
Description=Run SLZ platform database backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=slz-platform-backup.service

[Install]
WantedBy=timers.target