# packaging/systemd/slz-aiops.service.template
[Unit]
Description=SLZ AIOps Platform
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__ROOT_DIR__
ExecStart=__ROOT_DIR__/scripts/up.sh
ExecStop=__ROOT_DIR__/scripts/down.sh
RemainAfterExit=yes
TimeoutStartSec=600
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target